# conf/supervisor.conf
[program:telegrambots]
user=root
command=python /srv/www/telegrambots/bot.py
directory=/srv/www/telegrambots
autostart=true
autorestart=true
